# profile/profile.d/BUNDLE_GEM__FURY__IO
#!/bin/bash
lpass show "Gemfury" --field=api_token
